# s6/services.d/dev-server/run
#!/usr/bin/with-contenv sh


cd /app
while true; do
	exec python server.py
done
