# s6/services.d/uwsgi/run
#!/usr/bin/with-contenv sh

cd /app
exec uwsgi --ini wsgi.ini --plugin python3
